# deploy/systemd/mctrader-collector.service
[Unit]
Description=mctrader market data collector
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=120

[Service]
Type=simple
User=mctrader
Group=mctrader
WorkingDirectory=/opt/mctrader
ExecStart=/opt/mctrader/venv/bin/mctrader-collector
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mctrader-collector
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=MCTRADER_CONFIG_DIR=/opt/mctrader/config
EnvironmentFile=-/opt/mctrader/.env
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
