# deploy/launchd/com.mctrader.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mctrader.collector</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/mctrader/mctrader/venv/bin/mctrader-collector</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mctrader/mctrader</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>PYTHONDONTWRITEBYTECODE</key>
    <string>1</string>
    <key>MCTRADER_CONFIG_DIR</key>
    <string>/Users/mctrader/mctrader/config</string>
    <key>PATH</key>
    <string>/Users/mctrader/mctrader/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>KeepAlive</key>
  <dict>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/mctrader/Library/Logs/mctrader/collector.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mctrader/Library/Logs/mctrader/collector.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
